# ritooo.conf
[program:ritooo]
command=python3 /home/curator4/ritooo/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/ritooo/test.err.log
stdout_logfile=/home/curator4/ritooo/test.out.log